# docker/cluster-f-devops/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

# ── SysML kernel Postgres (loopback :5432) ───────────────────
# Must come up before the kernel; ephemeral data (create-drop on kernel boot).
[program:sysml-db]
command=/usr/local/bin/forge-sysml-db.sh
autostart=true
autorestart=true
startretries=3
startsecs=10
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── SysML v2 Java backend (Play dist, internal :8003) ────────
# Waits for sysml-db, then launches the staged Play app on SYSML_JAVA_PORT.
[program:sysml-java]
command=/usr/local/bin/forge-sysml-kernel.sh
directory=/opt/sysml
environment=JAVA_OPTS="%(ENV_JAVA_OPTS)s"
autostart=true
autorestart=true
startretries=3
startsecs=30
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── SysML FastAPI sidecar (public :9000 = SYSML_API_PORT) ────
[program:sysml-api]
command=/opt/sidecar-venv/bin/python -m uvicorn api_wrapper:app --host 0.0.0.0 --port %(ENV_SYSML_API_PORT)s --workers 1
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Forgejo (internal :3000, data on EFS /forge/repos) ───────
[program:forgejo]
command=/usr/local/bin/forgejo web --config /etc/forgejo/app.ini
directory=/forge/repos
user=git
environment=HOME="/var/lib/forgejo",GITEA_WORK_DIR="/forge/repos"
autostart=true
autorestart=true
startretries=5
startsecs=15
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── MinIO (S3 API :9090, console :9091, data on EFS /forge/minio)
# Runs as the uid-1000 "git" identity to match the /forge/minio EFS access
# point (posixUser 1000:1000). minio runs fine as a non-root uid.
[program:minio]
command=/usr/local/bin/minio server /forge/minio --address :%(ENV_MINIO_API_PORT)s --console-address :%(ENV_MINIO_CONSOLE_PORT)s
directory=/forge/minio
user=git
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Nginx front (public :80, /healthz + routing) ─────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
